# reboot.cron
* * * * /home/pi/projects/car/reboot.sh
